# doorbot-sync.service
[Unit]
Description=Sync doorbot repo from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=doorbot
ExecStart=/home/doorbot/doorbot/doorbot-sync.sh
